# supervisor.ini
[program:gin-hong-api]
directory=/golang/go/bin
command=/golang/go/bin/xueqiuv1.1 -env pro
autostart=true
autorestart=true
startsecs=1
user = root
stderr_logfile=/tmp/xueqiuv1.log
stdout_logfile=/tmp/xueqiuv1.log
redirect_stderr = true
stdout_logfile_maxbytes = 100MB
stdout_logfile_backups = 20


[program:gin-hong-api-xueqiuUgcV2]
directory=/golang/go/bin
command=/golang/go/bin/xueqiuUgcV2 -env pro
autostart=true
autorestart=true
startsecs=1
user = root
stderr_logfile=/tmp/xueqiuV2.log
stdout_logfile=/tmp/xueqiuV2.log
redirect_stderr = true
stdout_logfile_maxbytes = 100MB
stdout_logfile_backups = 20

[program:gin-hong-api-xueqiuCombinationUpdate]
directory=/golang/go/bin
command=/golang/go/bin/xueqiuCombinationUpdate -env pro
autostart=true
autorestart=true
startsecs=1
user = root
stderr_logfile=/tmp/xueqiuCombinationUpdate.log
stdout_logfile=/tmp/xueqiuCombinationUpdate.log
redirect_stderr = true
stdout_logfile_maxbytes = 100MB
stdout_logfile_backups = 20

[program:gin-hong-api-xueqiuCombination]
directory=/golang/go/bin
command=/golang/go/bin/xueqiuCombination -env pro
autostart=true
autorestart=true
startsecs=1
user = root
stderr_logfile=/tmp/xueqiuCombination.log
stdout_logfile=/tmp/xueqiuCombination.log
redirect_stderr = true
stdout_logfile_maxbytes = 100MB
stdout_logfile_backups = 20

[program:gin-hong-api-rank]
directory=/golang/go/bin
command=/golang/go/bin/pushRank -env pro
autostart=true
autorestart=true
startsecs=1
user = root
stderr_logfile=/tmp/pushRank.log
stdout_logfile=/tmp/pushRank.log
redirect_stderr = true
stdout_logfile_maxbytes = 100MB
stdout_logfile_backups = 20

[program:gin-hong-api-queue]
directory=/golang/go/bin
command=/golang/go/bin/queue -env pro
autostart=true
autorestart=true
startsecs=1
user = root
stderr_logfile=/tmp/queue.log
stdout_logfile=/tmp/queue.log
redirect_stderr = true
stdout_logfile_maxbytes = 100MB
stdout_logfile_backups = 20

[program:gin-hong-api-heshui]
directory=/golang/go/bin
command=/golang/go/bin/heshui -env pro
autostart=true
autorestart=true
startsecs=1
user = root
stderr_logfile=/tmp/queue.log
stdout_logfile=/tmp/queue.log
redirect_stderr = true
stdout_logfile_maxbytes = 100MB
stdout_logfile_backups = 20